# userspace/files_pixel3/sound.service
[Unit]
Description=Sound

[Service]
Type=oneshot
ExecStart=/usr/comma/sound/sound_init_pixel3.sh

[Install]
WantedBy=multi-user.target
